# rootfs/etc/services.d/ssh-mcp/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the SSH MCP Server
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare -a args

args+=("--host=$(bashio::config 'host')")
args+=("--port=$(bashio::config 'port')")
args+=("--user=$(bashio::config 'user')")

if bashio::config.has_value 'password'; then
  args+=("--password=$(bashio::config 'password')")
fi

if bashio::config.has_value 'ssh_key_path'; then
  args+=("--key=/config/$(bashio::config 'ssh_key_path')")
fi

if bashio::config.has_value 'sudo_password'; then
  args+=("--sudoPassword=$(bashio::config 'sudo_password')")
fi

if bashio::config.has_value 'su_password'; then
  args+=("--suPassword=$(bashio::config 'su_password')")
fi

if bashio::config.true 'disable_sudo'; then
  args+=("--disableSudo")
fi

args+=("--timeout=$(bashio::config 'timeout')")
args+=("--maxChars=$(bashio::config 'max_chars')")

if bashio::config.has_value 'api_key'; then
  # Operator pinned an explicit secret; use it as-is.
  args+=("--apiKey=$(bashio::config 'api_key')")
else
  # No pinned secret: auto-generate one (128-bit, URL-safe) and persist it in
  # /data, which Supervisor always mounts as this app's private persistent
  # storage — so the secret (and the MCP URL) survives restarts/updates.
  args+=("--secretPathFile=/data/secret_path.txt")
fi

args+=("--transport=http")
args+=("--httpPort=3000")

# The container always listens on 3000 internally; if the operator remapped it to a
# different host port via this app's Network tab, tell the server so its startup log
# prints a URL with the port that's actually reachable, not the internal one.
public_port="$(bashio::addon.port 3000)"
if bashio::var.has_value "${public_port}"; then
  args+=("--publicPort=${public_port}")
fi

bashio::log.info "Starting SSH MCP Server for $(bashio::config 'user')@$(bashio::config 'host'):$(bashio::config 'port') ..."

cd /app || bashio::exit.nok "Could not change to /app"
exec node build/index.js "${args[@]}"
